# etc/services.d/uwsgi/run
#!/usr/bin/execlineb -P
with-contenv
/usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi.conf -b 32768 --post-buffering=1